# cron/cron.hourly
#!/bin/sh

# Fetch data
if [ -x REPLACEINSTALLFOLDER/bin/plugins/REPLACEFOLDERNAME/fetch.pl ]
then
  cd REPLACEINSTALLFOLDER/bin/plugins/REPLACEFOLDERNAME/
  /usr/bin/perl fetch.pl -step 3600
fi

# Data import scheduler - Clean up hanging

jobnr=$(ls -d REPLACEINSTALLFOLDER/data/plugins/REPLACEFOLDERNAME/import/*.running.* | wc -l)
	
if [ $jobnr -ne 0 ]
	then 
		perl REPLACEINSTALLFOLDER/bin/plugins/REPLACEFOLDERNAME/scheduler.pl
fi
